# packages/panasonic-programmable-keys/panasonic_programmable_keys-0.2.4.tar.gz/panasonic_programmable_keys-0.2.4/src/panasonic_programmable_keys/install/templates/panasonic-programmable-keys-server.service.j2
[Unit]
Description=Start the Panasonic Keys system service, reading events off of the libinput device and sending them to the shared socket

[Service]
ExecStart={{ ppk_full_path }} server

[Install]
WantedBy=multi-user.target
