# packages/panasonic-programmable-keys/panasonic_programmable_keys-0.2.1.tar.gz/panasonic_programmable_keys-0.2.1/src/panasonic_programmable_keys/install/templates/panasonic-programmable-keys-client.service.j2
[Unit]
Description=Start the Panasonic Keys user service, reading events off of the server

[Service]
ExecStart={{ ppk_full_path }} client

[Install]
WantedBy=basic.target
